# systemd/shibui.service
[Unit]
Description=ShibUI Wayland Compositor
Documentation=https://github.com/yourusername/shibui
After=systemd-user-sessions.service getty.service
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/local/bin/shibui
Environment=WAYLAND_DISPLAY=wayland-1
Restart=on-failure
RestartSec=3
TimeoutStopSec=10

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
MemoryDenyWriteExecute=true

# Allow necessary permissions
ReadWritePaths=%h/.config/flowwm
ReadWritePaths=%h/.cache/flowwm

[Install]
WantedBy=graphical-session.target
